# skills/agent-chat/templates/systemd/everclaw-agent-chat.service
[Unit]
Description=EverClaw XMTP Agent Chat Daemon
Documentation=https://github.com/EverClaw/EverClaw
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory={{EVERCLAW_PATH}}/skills/agent-chat
ExecStart={{NODE_BIN}} {{EVERCLAW_PATH}}/skills/agent-chat/daemon.mjs
Restart=on-failure
RestartSec=10
Environment=NODE_ENV=production
Environment=HOME={{HOME}}

# Logging
StandardOutput=append:{{LOG_DIR}}/agent-chat.log
StandardError=append:{{LOG_DIR}}/agent-chat.err

# Security (user-level, no special privileges needed)
NoNewPrivileges=true

[Install]
WantedBy=default.target