# packages/core/templates/systemd/everclaw-guardian.service
[Unit]
Description=EverClaw Gateway Guardian (2-minute health check)
After=network.target morpheus-proxy.service

[Service]
Type=oneshot
ExecStart=__GUARDIAN_SCRIPT_PATH__
WorkingDirectory=__OPENCLAW_DIR__/workspace/scripts
Environment="HOME=__HOME__"
Environment="OPENCLAW_DIR=__OPENCLAW_DIR__"

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true

[Install]
WantedBy=default.target