# packages/core/templates/systemd/everclaw-guardian.timer
[Unit]
Description=EverClaw Gateway Guardian Timer (2-minute cycle)

[Timer]
OnBootSec=1min
OnUnitActiveSec=2min
Unit=everclaw-guardian.service

[Install]
WantedBy=timers.target