# packages/core/templates/systemd/morpheus-router.service
[Unit]
Description=Morpheus Proxy Router (Go binary)
After=network.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=__MORPHEUS_DIR__/proxy-router
WorkingDirectory=__MORPHEUS_DIR__
Environment="HOME=__HOME__"
Environment="MORPHEUS_DIR=__MORPHEUS_DIR__"

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
RestrictSUIDSGID=true

# Restart on failure
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target