# packages/core/templates/systemd/morpheus-proxy.service
[Unit]
Description=Morpheus OpenAI-Compatible Proxy (Node.js)
After=network.target morpheus-router.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=__NODE_PATH__ __PROXY_DIR__/morpheus-proxy.mjs
WorkingDirectory=__PROXY_DIR__
Environment="HOME=__HOME__"
Environment="MORPHEUS_DIR=__MORPHEUS_DIR__"

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
RestrictSUIDSGID=true

# Restart on failure
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target